# crontabs/nobody
*/15 * * * * /usr/local/bin/twitter_collector
