# crontabs/nobody
0 * * * * /usr/local/bin/wrap_cron /usr/local/bin/twitter_collector
0 * * * * /usr/local/bin/wrap_cron /usr/local/bin/github_collector
