# crontabs/nobody
0 * * * * /usr/local/bin/wrap_cron /usr/local/bin/twitter_collector rancher_labs ibuildthecloud
0 * * * * /usr/local/bin/wrap_cron /usr/local/bin/github_collector
0 * * * * /usr/local/bin/wrap_cron /usr/local/bin/github_release_collector 'rancherio/os' 'rancherio/python-agent' 'rancherio/host-api'
